# src/upstart/ceph-osd-all-reloader.conf
description "Ceph OSD (reload all instances)"

start on ceph-osd-all-reload or ceph-all-reload

task

script
    # someday upstart will let us reload all instances
    # with a wildcard or something similar; see
    # https://bugs.launchpad.net/upstart/+bug/1012938
    # until then,
    killall -1 ceph-osd
end script
